# rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.
cat < /dev/ttyACM0 > /tmp/heatingled

exit 0
